# service/30-apache2/run
#!/bin/sh
. /etc/service/template

msglog green "Starting Apache..."
exec apache2-foreground
