# service/90-glpi_init/run
#!/bin/sh
. /etc/service/template

msglog green "Starting glpi_init.sh..."
/root/glpi_init.sh
sleep infinity
